# supervisor-ghost.conf
[program:ghost]
priority=20
directory=/data/ghost
environment = NODE_ENV=production
command = npm start
user=root
autostart=true
autorestart=true
stdout_logfile=syslog
stderr_logfile=syslog
